# docker/supervisord.conf
; Sentrix Sail supervisor config.
;
; Extends the stock Laravel Sail supervisord (which only runs `php artisan
; serve`) with the two long-running realtime services so they start with the
; container and auto-restart — no more manual `reverb:start` / `horizon` after
; every `sail up` / `sail restart`.
;
; Mounted over /etc/supervisor/conf.d/supervisord.conf via compose.override.yml.
; Reverb runs INSIDE this container so the app's own broadcast push reaches it on
; localhost:8080 (same REVERB_HOST the browser uses via the published port).

[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

; Control socket so `supervisorctl status / restart reverb` works inside the
; container (stock Sail omits these, which is why supervisorctl had no socket).
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:php]
command=%(ENV_SUPERVISOR_PHP_COMMAND)s
user=%(ENV_SUPERVISOR_PHP_USER)s
environment=LARAVEL_SAIL="1"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:reverb]
command=/usr/bin/php /var/www/html/artisan reverb:start --host=0.0.0.0 --port=8080
user=%(ENV_SUPERVISOR_PHP_USER)s
environment=LARAVEL_SAIL="1"
autostart=true
autorestart=true
stopwaitsecs=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:horizon]
command=/usr/bin/php /var/www/html/artisan horizon
user=%(ENV_SUPERVISOR_PHP_USER)s
environment=LARAVEL_SAIL="1"
autostart=true
autorestart=true
; Give Horizon time to finish in-flight jobs on restart (graceful terminate).
stopwaitsecs=30
stopsignal=TERM
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
